# ish/ishbg
#!/sbin/openrc-run

description="iSH continues to run in the background"

command="/bin/cat"
command_args="/dev/location > /dev/null"
command_background="YES"

pidfile="/run/ishbg.pid"
